# src/polkit/burocracy/helper.service.in
# helper.service.in

[Unit]
Description=BeeKeeper privileged helper (DBus-activated)
After=dbus.service local-fs.target
Wants=local-fs.target

[Service]
Type=dbus
BusName=org.beekeeper.Helper
ExecStart=@BEEKEEPER_INSTALL_LIBDIR@/bin/beekeeper-helper
Restart=on-failure

# Hardening (keep, can extend later):
ProtectSystem=full
ReadWritePaths=/etc/bees
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
