# src/polkit/org.beekeeper.helper.service.in
[Unit]
Description=BeeKeeper privileged helper (DBus-activated)
After=dbus.service

[Service]
Type=dbus
BusName=org.beekeeper.Helper
ExecStart=@libexecdir@/beekeeper/beekeeper-helper
Restart=on-failure
# Harden:
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
# You can add more sandboxing depending on needed syscalls.

[Install]
WantedBy=multi-user.target
